# startup.service
[Unit]
Description=Install and run Dripper

[Service]
ExecStart=/opt/startup.sh

[Install]
WantedBy=multi-user.target
